# etc/profile.d/99-linuxmuster-linuxclient7.sh
scriptDir=$(/usr/sbin/linuxmuster-linuxclient7 get-constant scriptDir)
source $scriptDir/executeHookWithEnvFix.sh onLogin