# etc/systemd/system/sssd.service.d/override.conf
[Service]
TimeoutStartSec=10s
[Unit]
Before=systemd-user-sessions.service nss-user-lookup.target lightdm.service
